# templates/default/neo4j.init.erb
#!/usr/bin/env bash
#
# neo4j-<%= @instance_name %>
#
# chkconfig:   - 58 48
# description: neo4j-<%= @instance_name %>
# processname: neo4j-<%= @instance_name %>

# Source networking configuration
if [ -f /etc/sysconfig/network ]; then source /etc/sysconfig/network; fi

# Exit if networking is not up
[ "$NETWORKING" = "no" ] && exit

start() {
    # neo4j wrapper is smart enough to  change user and do similar things.
    sudo <%= @install_dir %>/bin/neo4j start
    return $?
}

stop() {
    sudo <%= @install_dir %>/bin/neo4j stop
    return $?
}

restart() {
    stop
    start
}

status() {
    sudo <%= @install_dir %>/bin/neo4j status
    return $?
}


case "$1" in 
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  status)
        status $2
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|status [-v]|}"
        exit 1
esac

exit $?
